# .config/systemd/user/commit_notes.service
[Unit]
Description=Run commit notes

[Service]
ExecStart=/home/bigforcegun/bin/commit_notes