# .config/systemd/user/gufw_icon.service
[Unit]
Description=GUFW tray icon

[Service]
Type=simple
ExecStart=/usr/bin/gufw_icon.sh

[Install]
WantedBy=default.target
