# Library/LaunchAgents/com.bigforcegun.mcp-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bigforcegun.mcp-proxy</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/com.bigforcegun.mcp-proxy.err.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/com.bigforcegun.mcp-proxy.out.log</string>
    <key>ProgramArguments</key>
    <array>
       <string>/bin/zsh</string>
       <string>-lc</string>
       <string>set -a; . ~/.config/mcp-proxy/env; exec mcp-proxy --config ~/.config/mcp-proxy/config.json</string>
    </array>
</dict>
</plist>
