# etc/systemd/system/mnt-station-backups.automount
[Unit]
Description=CIFS Station mount

[Automount]
Where=/mnt/station/backups
TimeoutIdleSec=0

[Install]
WantedBy=graphical.target